# home/dot_config/containers/systemd/piler.container
# piler — email archive with full-text search.
# Receives a BCC copy of every inbound + outbound from Postfix;
# indefinite retention (D-106, ~100GB for decades of 22-person email).
# UI at archive.sodimo.eu via CF Access.

[Unit]
Description=Piler email archive
Documentation=https://www.mailpiler.org/docs/
Wants=sodimo-network.service piler-mysql.service
After=sodimo-network.service piler-mysql.service

[Container]
Image=docker.io/sutoj/piler:1.4.5@sha256:fb3372fc01ac9f198589df43e8bb5b2a75589b8c297688ebefbcd03bd4203190
ContainerName=piler
Network=sodimo.network
PublishPort=127.0.0.1:8091:80
# SMTP sink (piler listens on 25 internally, we map externally to 2525
# so only Postfix can talk to it inside the podman network).
# piler-mysql.volume is mounted on piler-mysql (the MariaDB sidecar), NOT
# here. Piler reaches the server via MYSQL_HOSTNAME=piler-mysql in piler.env.
Volume=piler-data.volume:/var/piler
EnvironmentFile=%h/.config/containers/systemd/piler.env

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=sodimo.target
