# home/dot_config/containers/systemd/postfix.container
# postfix — outbound SMTP + inbound accepting for sodimo.eu + aliases.
# Vendored-upstream vendor class (Principle 1 counter-example): decades-old
# OSS, rock-solid, near-zero breakage surface at pinned version.
#
# Inbound: MX record points at harness.sodimo.eu; Postfix receives on 25.
# Outbound: email_send MCP tool on Cloudflare Worker writes to CF Queue
# `email_outbox`; sodimo-email-drain.service pulls and pipes to `sendmail -t`
# locally (Principle 3 pull-based wiring).

[Unit]
Description=Postfix MTA
Documentation=http://www.postfix.org/documentation.html
Wants=sodimo-network.service
After=sodimo-network.service

[Container]
Image=docker.io/boky/postfix:4.3.0
ContainerName=postfix
Network=sodimo.network
# Inbound SMTP (25), submission (587), SMTPS (465)
PublishPort=25:25
PublishPort=587:587
PublishPort=465:465
# main.cf + master.cf authored per sodimo/dotfiles#18 (see postfix-main.cf
# comment block for the design). boky/postfix merges POSTFIX_* env vars on
# top of the mounted main.cf at entrypoint time, so runtime overrides still
# win for hostname / TLS / milter tunables.
Volume=%h/.config/containers/systemd/postfix-main.cf:/etc/postfix/main.cf:ro,Z
Volume=%h/.config/containers/systemd/postfix-master.cf:/etc/postfix/master.cf:ro,Z
Volume=postfix-spool.volume:/var/spool/postfix
Volume=postfix-data.volume:/var/lib/postfix
EnvironmentFile=%h/.config/containers/systemd/postfix.env

[Service]
Restart=always

[Install]
WantedBy=sodimo.target
