# home/dot_config/containers/systemd/postfix.container
# postfix — outbound SMTP + inbound accepting for sodimo.eu + aliases.
# Vendored-upstream vendor class (Principle 1 counter-example): decades-old
# OSS, rock-solid, near-zero breakage surface at pinned version.
#
# Inbound: MX record points at harness.sodimo.eu; Postfix receives on 25.
# Outbound: email_send MCP tool on Cloudflare Worker writes to CF Queue
# `email_outbox`; sodimo-email-drain.service pulls and pipes to `sendmail -t`
# locally (Principle 3 pull-based wiring).

[Unit]
Description=Postfix MTA
Documentation=http://www.postfix.org/documentation.html
Wants=sodimo-network.service
After=sodimo-network.service

[Container]
Image=docker.io/boky/postfix:4.3.0
ContainerName=postfix
Network=sodimo.network
# Inbound SMTP (25), submission (587), SMTPS (465)
PublishPort=25:25
PublishPort=587:587
PublishPort=465:465
# main.cf + master.cf config mounts intentionally disabled — the files don't
# yet exist in the repo. Postfix uses env-based config (POSTFIX_* vars in
# postfix.env) which is sufficient for the outbound-via-CF-Queue posture.
# Real postfix config authoring is tracked at
# https://github.com/sodimo/dotfiles/issues/18; re-enable these mounts alongside
# that PR.
# Volume=%h/.config/containers/systemd/postfix-main.cf:/etc/postfix/main.cf:ro,Z
# Volume=%h/.config/containers/systemd/postfix-master.cf:/etc/postfix/master.cf:ro,Z
Volume=postfix-spool.volume:/var/spool/postfix
Volume=postfix-data.volume:/var/lib/postfix
EnvironmentFile=%h/.config/containers/systemd/postfix.env

[Service]
Restart=always

[Install]
WantedBy=sodimo.target
