# home/dot_config/containers/systemd/rspamd.container
# rspamd — spam filter + DKIM signing + reputation scoring.
# Postfix connects via milter on port 11332.

[Unit]
Description=rspamd
Documentation=https://rspamd.com/doc/
Wants=sodimo-network.service
After=sodimo-network.service

[Container]
Image=docker.io/rspamd/rspamd:3.9.1
ContainerName=rspamd
Network=sodimo.network
# Web UI (11334), milter (11332), controller (11333)
Volume=rspamd-data.volume:/var/lib/rspamd
Volume=rspamd-config.volume:/etc/rspamd/local.d

[Service]
Restart=always

[Install]
WantedBy=sodimo.target
