# home/dot_config/containers/systemd/vaultwarden.container
# vaultwarden — Bitwarden-compatible password vault, self-hosted.
# Sodimo canonical vault per `37-vault.md`. Employees log in via
# Bitwarden clients against vault.sodimo.eu (gated by CF Access).
# Break-glass via Tailscale only.

[Unit]
Description=Vaultwarden
Documentation=https://github.com/dani-garcia/vaultwarden/wiki
Wants=sodimo-network.service
After=sodimo-network.service

[Container]
Image=docker.io/vaultwarden/server:1.32.7@sha256:7a0aa23c0947be3582898deb5170ea4359493ed9a76af2badf60a7eb45ac36af
ContainerName=vaultwarden
Network=sodimo.network
# 2026-05-05: published on LAN + Tailscale (no 127.0.0.1 bind) so office
# devices reach http://192.168.0.134:8222 during demo. REVERT to
# `127.0.0.1:8222:80` once cloudflared + vault.sodimo.eu go live and the
# only path in is via Caddy on the host loopback.
PublishPort=8222:80
EnvironmentFile=%h/.config/sodimo/vaultwarden.env
# TM-01 (#27): vaultwarden-data is a named volume today; will move to a
# bind-mount under /var/lib/sodimo-data/vaultwarden/ when the bootc backup
# boundary lands.
Volume=vaultwarden-data.volume:/data

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=sodimo.target
